# hw3/hw3_3/webserver.service
[Unit]
Description=Simple bash HTTP server on port 80
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/simple-web-server.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
